# com.mlx.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mlx.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>REPO_DIR_PLACEHOLDER/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>REPO_DIR_PLACEHOLDER</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>REPO_DIR_PLACEHOLDER/logs/mlx-server.log</string>

    <key>StandardErrorPath</key>
    <string>REPO_DIR_PLACEHOLDER/logs/mlx-server.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
